# deploy/remy-bot.service
# systemd unit for the remy Telegram bot.
# Installed by deploy/bootstrap.sh to /etc/systemd/system/remy-bot.service
#
#   sudo systemctl enable --now remy-bot
#   journalctl -u remy-bot -f

[Unit]
Description=remy Telegram bot (auto TikTok editor)
Documentation=https://github.com/vyahello/remy
# the local Bot API container should be up first (docker compose)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=remy
Group=remy
WorkingDirectory=/opt/remy
# secrets live outside the repo — see deploy/env.example
EnvironmentFile=/etc/remy/env
ExecStart=/opt/remy/venv/bin/remy-bot
Restart=on-failure
RestartSec=5

# Good neighbor on a shared box: x265 encodes (ffmpeg children of this
# service) yield CPU to other services and are memory-capped inside this
# cgroup — when a render hits the cap, the kernel kills the ffmpeg (the
# bot catches it and reports a failed take) instead of dragging down
# whatever else runs on the host. Deliberately NO OOMScoreAdjust: the
# children would inherit it and shift kill pressure onto other services.
Nice=10
CPUWeight=50
IOWeight=50
MemoryHigh=60%
MemoryMax=75%
# Cap swap for this cgroup. Without it a render that exceeds MemoryMax
# spills gigabytes into swap and thrashes for 30+ min — making no
# progress AND starving the host's other services — instead of being
# killed. With it, an over-budget render OOMs fast and the bot reports
# a failed take. (A many-segment 1080p60 source can blow the budget;
# see the render notes about input fan-out.)
MemorySwapMax=256M
# hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/remy /home/remy /var/lib/telegram-bot-api

[Install]
WantedBy=multi-user.target
