# deploy/tokcut-bot.service
# systemd unit for the tokcut Telegram bot.
# Installed by deploy/bootstrap.sh to /etc/systemd/system/tokcut-bot.service
#
#   sudo systemctl enable --now tokcut-bot
#   journalctl -u tokcut-bot -f

[Unit]
Description=tokcut Telegram bot (auto TikTok editor)
Documentation=https://github.com/vyahello/tokcut
# the local Bot API container should be up first (docker compose)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=tokcut
Group=tokcut
WorkingDirectory=/opt/tokcut
# secrets live outside the repo — see deploy/env.example
EnvironmentFile=/etc/tokcut/env
ExecStart=/opt/tokcut/venv/bin/tokcut-bot
Restart=on-failure
RestartSec=5

# Good neighbor on a shared box: x265 encodes (ffmpeg children of this
# service) yield CPU to other services and are memory-capped inside this
# cgroup — when a render hits the cap, the kernel kills the ffmpeg (the
# bot catches it and reports a failed take) instead of dragging down
# whatever else runs on the host. Deliberately NO OOMScoreAdjust: the
# children would inherit it and shift kill pressure onto other services.
Nice=10
CPUWeight=50
IOWeight=50
MemoryHigh=60%
MemoryMax=75%
# hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/tokcut /home/tokcut /var/lib/telegram-bot-api

[Install]
WantedBy=multi-user.target
